# tests/hadoop/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:namenode]
command=/var/lib/hadoop/bin/hdfs namenode
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:datanode]
command=/var/lib/hadoop/bin/hdfs datanode
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:resourcemanager]
command=/var/lib/hadoop/bin/yarn resourcemanager
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nodemanager]
command=/var/lib/hadoop/bin/yarn nodemanager
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:krb5kdc]
command=krb5kdc -n
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
